# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Next.js standalone server — runs from /app/web so server.js resolves .next/ and public/
[program:web]
command=node /app/web/apps/web/server.js
directory=/app/web
priority=20
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; WebSocket server — runs from /app where node_modules and packages/ live
[program:ws]
command=/bin/sh -c "node_modules/.bin/tsx apps/ws/src/index.ts"
directory=/app
priority=20
autostart=true
autorestart=true
environment=NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Telegram bot — same monorepo root, no inbound port needed
[program:bot]
command=/bin/sh -c "node_modules/.bin/tsx apps/bot/src/index.ts"
directory=/app
priority=20
autostart=true
autorestart=true
environment=NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
